# templates/systemd.service.template
[Unit]
Description=APM Debug Host MCP Server
After=network.target

[Service]
Type=simple
User={{USER}}
WorkingDirectory={{MCP_HOME}}
ExecStart={{NODE_PATH}} {{MCP_HOME}}/src/index.js
Restart=always
RestartSec=10
Environment="NODE_ENV=production"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths={{MCP_HOME}}/logs

[Install]
WantedBy=multi-user.target